# www/forgejo/pkg/forgejo.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/forgejo"
daemon_user="_forgejo"

. /etc/rc.d/rc.subr

rc_bg=YES
rc_reload=NO

rc_cmd $1
